# workspace/crontab/laradock.example
* * * * * laradock /usr/bin/php /var/www/ff schedule:run >> /dev/null 2>&1
