# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/servo.news-backend/artisan schedule:run >> /dev/null 2>&1
