# workspace/crontab/root
* * * * * php /var/www/laravel/artisan schedule:run >> /dev/null 2>&1