# workspace/crontab/suitecrm
* * * * * laradock cd /var/www/crm; php -f cron.php > /dev/null 2>&1
