# workspace/conf/etc/supervisor.d/ssh.conf
[group:ssh]
programs=sshd
priority=30

[program:sshd]
command = /opt/docker/bin/service.d/ssh.sh
process_name=%(program_name)s
startsecs=0
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
